# etc/exporterite confid/prometheus-mktxp-exporter.service
[Unit]
Description=MKTXP Exporter

[Service]
User=user
ExecStart=/home/user/.local/bin/mktxp export

[Install]
WantedBy=default.target